# modules/hosting/files/default.target
[Unit]
Description=Default user target
AllowIsolate=true
