# modules/hosting/templates/mono-fcgi.service.erb
[Unit]
Description=FastCGI mono process for <%= @destination %>

[Service]
WorkingDirectory=<%= @base_dir %>/apps/<%= @destination %>
# make socket accessible for admin group members
UMask=007
ExecStart=/usr/bin/fastcgi-mono-server4 /socket=unix:<%= @base_dir %>/run/<%= @destination %>.sock /appconfigdir=<%= @base_dir %>/etc/mono-<%= @destination %> /printlog=True /loglevels=Standard

# extended features
RestartSec=1
Restart=always

[Install]
WantedBy=default.target
