# modules/hosting/templates/php5-fpm.service.erb
[Unit]
Description=FastCGI PHP5 process for <%= @destination %>

[Service]
WorkingDirectory=<%= @base_dir %>/apps/<%= @destination %>
# make socket accessible for admin group members
UMask=007
ExecStartPre=/usr/sbin/php5-fpm --fpm-config <%= @base_dir %>/etc/php5-<%= @destination %>/php-fpm.conf --test
ExecStart=/usr/sbin/php5-fpm --fpm-config <%= @base_dir %>/etc/php5-<%= @destination %>/php-fpm.conf

PIDFile=<%= @base_dir %>/run/php5-fpm-<%= @destination %>.pid

# extended features
RestartSec=1
Restart=always

[Install]
WantedBy=default.target
